# fluxcd/environment/nebius-cloud/base/custom-configmaps-soperator/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null ; Output only to stdout/stderr
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:slurmd]
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=true
command=/opt/bin/slurm/slurmd_entrypoint.sh
autostart=true
autorestart=true
startsecs=0
stopasgroup=true ; Send SIGTERM to all child processes of supervisord
killasgroup=true ; Send SIGKILL to all child processes of supervisord
stopsignal=SIGTERM ; Signal to send to the program to stop it
stopwaitsecs=10 ; Wait for the process to stop before sending a SIGKILL

[program:sshd]
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=true
command=/usr/sbin/sshd -D -e -f /mnt/ssh-configs/sshd_config
autostart=true
autorestart=true
startsecs=0
stopasgroup=true ; Send SIGTERM to all child processes of supervisord
killasgroup=true ; Send SIGKILL to all child processes of supervisord
stopsignal=SIGTERM ; Signal to send to the program to stop it
stopwaitsecs=10 ; Wait for the process to stop before sending a SIGKILL

[program:dockerd]
priority=1
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=true
command=/usr/bin/dockerd
autostart=true
autorestart=true
startsecs=0
stopasgroup=true ; Send SIGTERM to all child processes of supervisord
killasgroup=true ; Send SIGKILL to all child processes of supervisord
stopsignal=SIGTERM ; Signal to send to the program to stop it
stopwaitsecs=10 ; Wait for the process to stop before sending a SIGKILL

[program:nv-hostengine]
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
redirect_stderr=true
command=/usr/bin/nv-hostengine --no-daemon --log-level ERROR --service-account nvidia-dcgm
autostart=true
autorestart=true
startsecs=0
stopasgroup=true ; Send SIGTERM to all child processes of supervisord
killasgroup=true ; Send SIGKILL to all child processes of supervisord
stopsignal=SIGTERM ; Signal to send to the program to stop it
stopwaitsecs=10 ; Wait for the process to stop before sending a SIGKILL
